# supervisor/supervisord_config/config.d/format_data.ini
[program:format_data]
directory=/opt/Alarm/format_data/
command=python3 /opt/Alarm/format_data/take_data.py 
user=root
autorestart=true

