# supervisor/supervisord_config/config.d/post_mail.ini
[program:mail_raw]
directory=/opt/Alarm/post_mail/
command=python3 /opt/Alarm/post_mail/mail_raw.py
user=root
autorestart=true
