# supervisor/supervisord_config/config.d/send_top.ini
[program:send_top]
directory=/opt/Alarm/send_top/
command=python3 /opt/Alarm/send_top/send_data.py 
user=root
autorestart=true

